# amdgpu-fancontrol.service
[Unit]
Description=amdgpu-fancontrol
StartLimitIntervalSec=10
StartLimitBurst=3

[Service]
Type=simple
ExecStartPre=/bin/sleep 1
ExecStart=/usr/bin/amdgpu-fancontrol
Restart=on-failure

[Install]
WantedBy=multi-user.target
