# services/amdgpu-fancontrol.path
[Path]
PathExistsGlob=/sys/class/drm/card*/device/hwmon/hwmon?/pwm1

[Install]
WantedBy=default.target
